# legacy/hosted/supervisord.conf
[supervisord]
nodaemon=true

[program:nuoagent]
user=nuodb
command=/opt/nuodb/etc/nuoagent start
